# archive/mailadmin/misc/cron.d/qstat
*/1 * * * * root        /home/matsu/public_html/mailadmin/cron/stat.sh > /dev/null
