# archive/mailadmin/misc/cron.d/send
*/1 * * * * matsu       /usr/local/bin/php -q  /home/matsu/public_html/mail/cron/SendScheduleAction.php > /dev/null
